# package/luci-theme-noobwrt/root/etc/uci-defaults/30_luci-theme-noobwrt
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci get luci.themes.NoobWrt >/dev/null 2>&1 || \
	uci batch <<-EOF
		set luci.themes.NoobWrt=/luci-static/noobwrt
		set luci.main.mediaurlbase=/luci-static/noobwrt
		commit luci
	EOF
fi

exit 0

