# package/luci-app-aw1k-led/files/etc/uci-defaults/99-set
#!/bin/sh

chmod +x /usr/bin/led-status.sh
chmod +x /etc/init.d/led-status

/etc/init.d/led-status enable
/etc/init.d/led-status start

exit 0
